# luci-app-timecontrol/root/etc/uci-defaults/luci-app-timecontrol
#!/bin/sh

uci -q batch <<-EOF01 >/dev/null
	delete firewall.timecontrol
	set firewall.timecontrol=include
	set firewall.timecontrol.type=script
	set firewall.timecontrol.path=/var/etc/timecontrol.include
	set firewall.timecontrol.reload=1
EOF01

uci -q batch <<-EOF02 >/dev/null
	delete ucitrack.@timecontrol[-1]
	add ucitrack timecontrol
	set ucitrack.@timecontrol[-1].init=timecontrol
	commit ucitrack
EOF02

# ★★★ 关键修复5：初始化配置文件 ★★★
if [ ! -f /etc/config/timecontrol ]; then
    touch /etc/config/timecontrol
    uci -q batch <<-EOF03 >/dev/null
        set timecontrol.config=basic
        set timecontrol.config.enabled='0'
        set timecontrol.config.scan_interval='10'
        commit timecontrol
EOF03
fi

rm -rf /tmp/luci-*cache
exit 0
